# Install/Files/readonly.service
[Unit]
Description=SPQR readonly.service

[Service]
Type=simple
ExecStart=/bin/bash -c 'sudo umount /'

[Install]
WantedBy=default.target
